# packaging/systemd/rillan.service
# SPDX-FileCopyrightText: 2026 Rillan AI LLC
# SPDX-License-Identifier: Apache-2.0

[Unit]
Description=Rillan local OpenAI-compatible proxy
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
ExecStart=%h/.local/bin/rillan serve --config %h/.config/rillan/config.yaml
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
